# systemd/aptly-api.service.d/10-no-lock.conf
[Service]
ExecStart=
ExecStart=/usr/bin/aptly api serve \
            -no-lock \
            -config=/etc/aptly-api.conf \
            -listen=${LISTEN_ADDRESS}
