# systemd/appstream-generate.timer
[Unit]
Description=Periodically generate appstream metadata

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
