# systemd/appstream-generate.service
[Unit]
Description=Generate appstream metadata

[Service]
Type=oneshot
ExecStart=/usr/bin/appstream-generate
User=aptly-api
Group=aptly-api
